# ansible/common/files/system/cluster.target
[Unit]
Description=Cluster Services
After=docker.service

[Install]
WantedBy=server.target