# ansible/common/files/system/rc-local.service
[Unit]
Description=Rc.local Script Trigger

[Service]
Type=oneshot
ExecStart=/usr/bin/rc-local

[Install]
WantedBy=server.target
